# macaw/docker/ic_app/config/gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn app:app -b localhost:5001
directory=/deploy/app
